# config/includes.chroot/etc/profile.d/aerthos.sh
# some security enhancements via env
export QT_ACCESSIBILITY=0
export NO_AT_BRIDGE=1
export TZ=UTC